# roles/rclone/files/rclone-storagebox.service
[Unit]
Description=RClone Backup

[Service]
Type=simple
ExecStart=/usr/local/bin/rclone-storagebox
